# templates/vaultwarden-backup.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Timer backup Vaultwarden quotidien

[Timer]
OnCalendar={{ vw_backup_calendar }}
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
